# deploy/t212.service
[Unit]
Description=T212 Portfolio Dashboard
After=network-online.target
Wants=network-online.target

[Service]
EnvironmentFile=/etc/t212/config.env
Environment="SIGNAL_CLI_CONFIG=/var/lib/t212/signal-cli"
ExecStart=/usr/local/bin/t212 serve
Restart=always
RestartSec=5
User=t212
Group=t212
NoNewPrivileges=yes
ProtectSystem=strict
ReadWritePaths=/var/lib/t212
PrivateTmp=yes
ProtectHome=yes
PrivateDevices=yes
CapabilityBoundingSet=

[Install]
WantedBy=multi-user.target
